# packages/obsidian-desktop/src/etc/xdg/autostart/obsidian-waybar.desktop
[Desktop Entry]
Type=Application
Name=Obsidian Waybar
Exec=waybar
OnlyShowIn=Obsidian;
X-GNOME-Autostart-enabled=true
