# packages/obsidian-welcome/src/etc/xdg/autostart/obsidian-welcome.desktop
[Desktop Entry]
Type=Application
Name=Obsidian Welcome
Exec=obsidian-welcome
OnlyShowIn=Obsidian;
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
Terminal=false
Categories=System;
